# codAPI.service.example
[Unit]
Description=Call Of Duty API
After=multi-user.target

[Service]
Type=simple
Restart=always
ExecStart=/usr/bin/python3 /home/<user>/callofdutyapi/main.py

[Install]
WantedBy=multi-user.target