# codAPI.service.example
[Unit]
Description=Call Of Duty API
After=multi-user.target

[Service]
Type=idle
ExecStart=/usr/bin/python3 /home/<USER>/callofdutyapi/main.py

[Install]
WantedBy=multi-user.target